# tg-tui-webapp-tunnel.service
[Unit]
Description=Cloudflare Tunnel for TG-TUI Webapp
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
User=aid
Group=aid
WorkingDirectory=/home/aid/workspace/tg-tui
ExecStart=/usr/local/bin/cloudflared tunnel --config /home/aid/workspace/tg-tui/cloudflared-config.yml run tg-tui-webapp
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target